# docker/supervisord.conf
[supervisord]
logfile=/dev/null
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
user=root

[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=php-fpm -R
autostart=true
autorestart=true
user=root
numprocs=1